# net/alfred/files/mesh-routing.config
#!/bin/sh /etc/rc.common

START=96
STOP=10

USE_PROCD=1

start_service() {
    # Enable cross-mesh client routing through batman-adv
    logger "mesh-routing: Setting up cross-mesh client routes"

    # Wait for network services to be ready
    sleep 45

    # Add routes for cross-mesh connectivity
    (
        # Wait for bat0 interface to be fully operational
        local retry=0
        while [ $retry -lt 30 ]; do
            if ip link show bat0 >/dev/null 2>&1 && ip addr show bat0; then
                logger "mesh-routing: bat0 interface is ready"
                break
            fi
            retry=$((retry + 1))
            logger "mesh-routing: Waiting for bat0 interface (attempt $retry/30)"
            sleep 2
        done

        # Enable IP forwarding for cross-mesh client communication
        echo 1 > /proc/sys/net/ipv4/ip_forward
        logger "mesh-routing: IP forwarding enabled"

        # Add routing rules for 10.41.0.0/16 mesh client subnet
        if ! ip route show | grep -q "10.41.0.0/16"; then
            ip route add 10.41.0.0/16 dev bat0 metric 100 2>/dev/null || true
            logger "mesh-routing: Added route for mesh client subnet"
        fi

        # Enable multicast forwarding for batman-adv
        echo 1 > /proc/sys/net/ipv4/conf/bat0/mc_forwarding 2>/dev/null || true
        echo 1 > /proc/sys/net/ipv4/conf/br-ahwlan/mc_forwarding 2>/dev/null || true

        # Setup ARP proxying for cross-mesh client discovery
        echo 1 > /proc/sys/net/ipv4/conf/bat0/proxy_arp 2>/dev/null || true
        echo 1 > /proc/sys/net/ipv4/conf/br-ahwlan/proxy_arp 2>/dev/null || true

        logger "mesh-routing: Cross-mesh client routing configured"
    ) &
}

stop() {
    logger "mesh-routing: Stopping mesh routing service"
}